# packages/server/transcription-api/scripts/transcription-api.service
[Unit]
Description=Transcription service
After=network.target

[Service]
User=transcribe
Group=transcribe
WorkingDirectory=/usr/local/transcribe
ExecStart=/usr/bin/env gunicorn --access-logfile - --error-logfile - --timeout 0 --workers 3 --bind=0.0.0.0:8014 routes.routes:app
Restart=always
EnvironmentFile=/etc/transcription/environment

[Install]
WantedBy=multi-user.target